# ansible/roles/instance.bitcoind_automated_updates/templates/automated-updates.service.j2
[Unit]
Description=bitcoin-node stop, update, and reboot service.

[Service]
Type=simple
Restart=no
ExecStart=/usr/local/sbin/update.sh
[Install]
WantedBy=multi-user.target
