# ansible/roles/instance.bitcoind_automated_updates/templates/automated-updates.timer.j2
[Unit]
Description=Runs the Bitcoin updater.

[Timer]
OnCalendar=daily
RandomizedDelaySec=21600
Unit=automated-updates.service

[Install]
WantedBy=multi-user.target
